# rust-agent/svf-print-agent.service
[Unit]
Description=SVF Print Agent — CUPS printing with MQTT orchestration
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
User=svf-agent
Group=lpadmin
ExecStart=/usr/local/bin/svf-print-agent
Environment=RUST_LOG=info,svf_print_agent=debug
WorkingDirectory=/var/lib/svf-agent
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/svf-agent /tmp/svf-agent
PrivateTmp=true

[Install]
WantedBy=multi-user.target
